# floginapp-infra/templates/default/floginapp-gunicorn.conf.erb
[program:<%= node['floginapp-infra']['app_name'] %>]
command = gunicorn wsgi:app -b #{node['floginapp-infra']['gunicorn_socket']}
directory = <%= node['floginapp-infra']['app_dir'] %>